# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
user=postgres
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/etc/postgresql/15/main/postgresql.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log

[program:officeplane-api]
user=appuser
directory=/app
command=/opt/venv/bin/uvicorn officeplane.api.main:app --host 0.0.0.0 --port 8001
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/officeplane-api.log
stderr_logfile=/var/log/supervisor/officeplane-api_err.log
environment=DATABASE_URL="postgresql://officeplane:officeplane@localhost:5432/officeplane",PATH="/opt/venv/bin:%(ENV_PATH)s"

[program:officeplane-ui]
user=appuser
directory=/app/ui
command=/usr/bin/node .next/standalone/ui/server.js
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/officeplane-ui.log
stderr_logfile=/var/log/supervisor/officeplane-ui_err.log
environment=DATABASE_URL="postgresql://officeplane:officeplane@localhost:5432/officeplane",NEXT_PUBLIC_API_URL="http://localhost:8001",NEXT_PUBLIC_WS_URL="ws://localhost:8001",PORT="3000",HOSTNAME="0.0.0.0"
